# caulklinux/iso/airootfs/etc/profile.d/installer.sh
#!/bin/sh
# auto-launch installer when root logs into tty1
if [ "$(id -u)" -eq 0 ] && [ "$(tty)" = "/dev/tty1" ]; then
    exec caulk-install
fi
